# home/.config/systemd/user/default.target
[Unit]
Description=My unprivileged services