# home/.config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon
After=dbus.service

[Service]
ExecStart=/usr/bin/dunst -print

[Install]
WantedBy=default.target
