# home/.config/systemd/user/unison.service
[Unit]
Description=Unison file synchronizer

[Service]
Type=oneshot
ExecStart=/usr/bin/unison sync -batch