# home/.config/systemd/user/unison.path
[Unit]
Description=Unison inotify activation
Before=paths.target

[Path]
PathChanged=%h/sync/unison

[Install]
WantedBy=default.target
